# confs/freeDNS.service
[Service]
ExecStart=/usr/bin/freeDNS
StandardOutput=journal
StandardError=journal
Restart=on-failure
